# raspi-wol@.service
[Unit]
Description=Raspberry Pi Wake-on-LAN Service for %i
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/threedpro/raspi-wol/raspi-wol.py %I
WorkingDirectory=/home/threedpro/raspi-wol/
Restart=always
PIDFile=/run/raspi-wol-%i.pid

[Install]
WantedBy=multi-user.target
